# foo.service
[Unit]
Description=Foo hello world

[Service]
ExecStart=/usr/bin/echo HELLO WORLD
