# foo.timer
[Unit]
Description=Run foo hourly

[Timer]
OnCalendar=hourly
 
[Install]
WantedBy=timers.target
