# getty@tty2.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin hendry --noclear %I 38400 linux
